# deploy/ec2/wtg-mci-api.service
# WTG mci-api — /v1/tx (매매 요청 → broker). flag 는 이전 docker-compose.prod.yml 과 동일.
[Unit]
Description=WTG mci-api — /v1/tx (매매 요청 → broker)
After=network-online.target wtg-etcd.service
Wants=network-online.target wtg-etcd.service
# broker/etcd 늦게 떠도 자동 복구되도록 start limit 해제 (무한 재시도)
StartLimitIntervalSec=0

[Service]
Type=simple
User=winway
Group=winway
WorkingDirectory=/home/winway/nh-fxallone-server/wtg
EnvironmentFile=/home/winway/nh-fxallone-server/wtg/wtg.env
# dev EC2 한정 — tx 전문 debug 로그 (wtg.env 는 CI 가 매 배포 덮어쓰므로 여기 고정).
# 운영 전환 시 이 줄 제거.
Environment=WTG_API_LOG_LEVEL=debug
# --etcd: user-profiles (usid→Site/Tier) watch + 라우팅 룰 동기화.
# --jwt-key: RS256 발급 — login 이 access_token 반환, /v1/refresh 활성.
#            키는 setup-ec2.sh 가 1회 생성 (jwt-private.pem / jwt-public.pem).
# --dev: dev EC2 한정 — admin UI DevMode (X-WTG-User) 의 tx-test 가 mci-api
#        관문을 통과하게 허용. Bearer 가 오면 JWT 검증이 우선 적용되므로
#        (middleware/auth.go authenticate) JWT 경로도 공존. 운영 전환 시 제거 필수.
# --login-mode=chain: NH 엔진 인증 사슬 (W1101S02→W1130A02, lgnIdntCon 보관).
#        svc-inc-dir 필수 — docs/auth.md §3.1.
# --login-skip-cert: ① 인증서 인증(W1101S02) 건너뛰고 data.lgnId 로 세션만 개설.
#        인증서 미적용 과도기 — 운영 전환 시 이 줄 제거 (사용자는 CSC004M 존재 필요).
ExecStart=/home/winway/nh-fxallone-server/wtg/bin/mci-api \
    --listen=:8080 \
    --broker-host=${BROKER_HOST} \
    --broker-port=${BROKER_PORT} \
    --dev \
    --dev-routes-file=/home/winway/nh-fxallone-server/wtg/etc/wtg-routes.json \
    --etcd=http://127.0.0.1:2379 \
    --user-profiles-prefix=wtg/auth/user-profiles/ \
    --jwt-key=/home/winway/nh-fxallone-server/wtg/etc/jwt-private.pem \
    --svc-common-header=/home/winway/nh-fxallone-server/win/src/inc/com/comhdr.h \
    --svc-inc-dir=/home/winway/nh-fxallone-server/win/src/inc/trn \
    --login-mode=chain \
    --login-skip-cert

Restart=on-failure
RestartSec=10
TimeoutStopSec=30
LimitNOFILE=65536

StandardOutput=journal
StandardError=journal
SyslogIdentifier=wtg-mci-api

[Install]
WantedBy=multi-user.target
