# deploy/ec2/wtg-mci-edge-md.service
# WTG mci-edge-md — FIX 4.4 시세 gateway. flag 는 이전 docker-compose.prod.yml 과 동일.
[Unit]
Description=WTG mci-edge-md — FIX 4.4 시세 gateway (5011) + stats (5012)
After=network-online.target wtg-etcd.service wtg-mci-price.service
Wants=network-online.target wtg-etcd.service wtg-mci-price.service
# broker/etcd 늦게 떠도 자동 복구되도록 start limit 해제 (무한 재시도)
StartLimitIntervalSec=0

[Service]
Type=simple
User=winway
Group=winway
WorkingDirectory=/home/winway/nh-fxallone-server/wtg
EnvironmentFile=/home/winway/nh-fxallone-server/wtg/wtg.env
ExecStart=/home/winway/nh-fxallone-server/wtg/bin/mci-edge-md \
    --port=5011 \
    --stats=127.0.0.1:5012 \
    --sender=WTG_MD \
    --etcd=http://127.0.0.1:2379 \
    --upstream=127.0.0.1:50051 \
    --seed-cp=ECN_MD_TEST_01=test-pw,HQ,VIP,ECN_MD_TEST_01

Restart=on-failure
RestartSec=10
TimeoutStopSec=30
LimitNOFILE=65536

StandardOutput=journal
StandardError=journal
SyslogIdentifier=wtg-mci-edge-md

[Install]
WantedBy=multi-user.target
