# deploy/ec2/wtg-fx-sync.service
# WTG fx-sync — 고객 등급(user_profile) + 고객 스프레드(customer_margin) 를 etcd
# 로 미러. oneshot — timer 가 주기 실행.
#
# ⚠️ --table 은 user_profile,customer_margin 만 (--table=all 아님):
#   base 마진(hq/site/swap)·통화·페어는 seed-catalog.sh 의 pricing.json 이 소유.
#   fx-sync 가 --table=all 로 그 레이어를 db-mirror 빈/샘플로 덮어쓰면 충돌.
#   Oracle backend 로 전 마스터 데이터를 fx-sync 가 소유하게 되면 --table=all +
#   seed-catalog 은퇴 (docs/customer-spread-sync.md / auth.md §6.5).
#
# 현재 source=file (etc/db-mirror/{user_profile,customer_margin}.json).
# Oracle 확정 시 --source=oracle + DSN 으로 교체.
[Unit]
Description=WTG fx-sync — 마스터 데이터 → etcd 미러 (통화/마진/tier/고객스프레드)
After=network-online.target wtg-etcd.service
Wants=network-online.target wtg-etcd.service

[Service]
Type=oneshot
User=winway
Group=winway
WorkingDirectory=/home/winway/nh-fxallone-server/wtg
EnvironmentFile=/home/winway/nh-fxallone-server/wtg/wtg.env
ExecStart=/home/winway/nh-fxallone-server/wtg/bin/fx-sync \
    --source=file \
    --source-dir=/home/winway/nh-fxallone-server/wtg/etc/db-mirror \
    --etcd=127.0.0.1:2379 \
    --table=user_profile,customer_margin

StandardOutput=journal
StandardError=journal
SyslogIdentifier=wtg-fx-sync
