# deploy/systemd/wtg-mci-edge-price.service
[Unit]
Description=WTG mci-edge-price — DMZ : 시세 ws 외부 노출
Documentation=file:///opt/wtg/share/doc/operations.md
After=network-online.target wtg-mci-price.service
Wants=network-online.target

[Service]
Type=simple
User=wtg
Group=wtg
EnvironmentFile=/etc/wtg/wtg.env
ExecStart=/opt/wtg/bin/mci-edge-price \
    --listen :8083 \
    --tls-cert /etc/pki/wtg/dmz.crt \
    --tls-key  /etc/pki/wtg/dmz.key \
    --upstream int1.internal:50051,int2.internal:50051 \
    --grpc-tls-ca   ${WTG_BROKER_TLS_CA} \
    --grpc-tls-cert ${WTG_BROKER_TLS_CERT} \
    --grpc-tls-key  ${WTG_BROKER_TLS_KEY} \
    --quote-stream \
    --customer-stream \
    --otel-endpoint ${WTG_OTEL_ENDPOINT}

Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60
TimeoutStopSec=30

LimitNOFILE=65536
LimitNPROC=4096

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log/wtg
ProtectHome=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=wtg-mci-edge-price

[Install]
WantedBy=multi-user.target
